# tools/ubuntu/etc/systemd/system/armcnc_display_off.service
[Unit]
Description=armcnc_display_off
After=graphical.target

[Service]
Type=simple
User=sunrise
WorkingDirectory=/opt/armcnc/
Environment=DISPLAY=:0.0
Environment=XAUTHORITY=/home/sunrise/.Xauthority
ExecStart=/usr/bin/xrandr --output DSI-1 --off
Restart=on-failure

[Install]
WantedBy=multi-user.target
